# scheduler/com.mab.daily-slicer-undeploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mab.daily-slicer-undeploy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mabdelouahab/mab_lab/mab_code/My_Daily_Slicer/scheduler/run-undeploy.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mabdelouahab/mab_lab/mab_code/My_Daily_Slicer/scheduler/logs/undeploy.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mabdelouahab/mab_lab/mab_code/My_Daily_Slicer/scheduler/logs/undeploy.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
